# contrib/systemd/steamdatabasebackend.service
[Unit]
Description=Steam Database Backend
Wants=network-online.target
After=network-online.target

[Service]
Type=simple
User=steamdb
Group=steamdb
WorkingDirectory=/opt/steamdatabasebackend
ExecStart=/opt/steamdatabasebackend/SteamDatabaseBackend
EnvironmentFile=-/etc/default/steamdatabasebackend
RuntimeDirectory=steamdatabasebackend
RuntimeDirectoryMode=0750
Restart=always
RestartSec=15
TimeoutStopSec=30
KillSignal=SIGTERM
StandardOutput=journal
StandardError=journal
SyslogIdentifier=steamdatabasebackend

[Install]
WantedBy=multi-user.target
